# experimental/container_deployments/aws_ec2_experiments/systemd_override_approach/network-overrides/etc/systemd/system/systemd-networkd.socket.d/10-disable.conf
[Unit]
Description=Disable networkd socket inside apptainer booted container
ConditionPathExists=!/etc/systemd/networkd.allow

[Socket]
ListenStream=
ListenDatagram=
ListenNetlink=
