# experimental/container_deployments/aws_ec2_experiments/systemd_override_approach/network-overrides/etc/systemd/system/systemd-networkd.service.d/10-disable.conf
[Unit]
Description=Disable networkd inside apptainer booted container
ConditionPathExists=!/etc/systemd/networkd.allow

[Service]
ExecStart=
ExecStart=/bin/false
