# experimental/container_deployments/aws_ec2_experiments/systemd_override_approach/network-overrides/etc/systemd/system/NetworkManager.service.d/10-disable.conf
[Unit]
Description=Disable NetworkManager inside apptainer booted container
ConditionPathExists=!/etc/NetworkManager/allow-start

[Service]
ExecStart=
ExecStart=/bin/false
