# experimental/container_deployments/aws_ec2_experiments/systemd_override_approach/network-overrides/etc/systemd/system/systemd-resolved.service.d/10-disable.conf
[Unit]
Description=Disable systemd-resolved inside apptainer booted container
ConditionPathExists=!/etc/systemd/resolved.allow

[Service]
ExecStart=
ExecStart=/bin/false
